# overlays/orbital/etc/init.d/docker-compose-autostart
#!/sbin/openrc-run

depend() {
	after docker net
}

start() {
	/usr/local/bin/compose-autostart.py /mnt/share/compose
}
